# picam-stream.service
# systemd unit TEMPLATE for picam-stream.
# install.sh substitutes __USER__ and __APPDIR__ and writes the result to
# /etc/systemd/system/picam-stream.service. To install by hand instead, copy this
# file there and replace the two placeholders with your username and the absolute
# path to this project directory.
[Unit]
Description=Picamera2 MJPEG camera stream (Camera Module 2)
After=network.target

[Service]
Type=simple
User=__USER__
ExecStart=/usr/bin/python3 __APPDIR__/mjpeg_server.py
Restart=on-failure
RestartSec=5
Nice=5

# Bind to localhost by default -- a reverse proxy (see Caddyfile.example) terminates
# TLS + auth. Set PICAM_BIND=0.0.0.0 here to serve directly on the LAN (unauthenticated).
Environment=PICAM_BIND=127.0.0.1
# Runtime state (stream-config.json) lives in /var/lib/picam-stream, owned by the
# service user. STATE_DIRECTORY is exported to the process.
StateDirectory=picam-stream
RuntimeDirectory=picam-stream

# --- sandboxing -----------------------------------------------------------------
# High-value hardening for a network-facing service. The camera needs a broad device
# surface (/dev/video*, /dev/media*, /dev/dma_heap*), so we deliberately do NOT use
# PrivateDevices / DevicePolicy=closed -- that would break capture. If the camera fails
# to start after enabling this, check `journalctl -u picam-stream -e` and relax the
# offending directive (most likely ProtectSystem or SystemCallFilter).
NoNewPrivileges=yes
CapabilityBoundingSet=
AmbientCapabilities=
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
ProtectProc=invisible
ProtectControlGroups=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectClock=yes
ProtectHostname=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
SupplementaryGroups=video
# MemoryDenyWriteExecute is left OFF: it can break the JPEG encoder's native libraries.
# Enable only after confirming the stream still works on real hardware.

[Install]
WantedBy=multi-user.target
